# scripts/bgp/etc/rc.local
#!/bin/bash
/opt/src/trigger_ipv46.sh &
exit 0
