# scripts/lamp/etc/rc.local
#!/bin/bash
/opt/src/trigger_hosts.sh &
exit 0
